# docker/metaxplor-hpc/runit/opal-tomcat/run
#!/bin/sh

/usr/local/tomcat/bin/replace_var_opal_prop.sh
exec sudo -u tomcat /usr/local/tomcat/bin/catalina.sh run